# projects/2605101/doc-viewer/doc-viewer.service
[Unit]
Description=Doc Viewer - Markdown/HTML Upload & Preview Service
After=network.target

[Service]
Type=simple
User=root
Environment=DOC_HOST=doc.20100706.xyz
Environment=DOC_PORT=80
Environment=DOC_DATA_DIR=/data/doc-viewer/data
Environment=DOC_RETENTION_DAYS=30
ExecStart=/data/doc-viewer/venv/bin/python /data/doc-viewer/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
